# meta-bunch-linux/recipes-core/bunch-update/bunch-update/update.timer
[Unit]
Description=Update timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=10

[Install]
WantedBy=timers.target