# meta-bunch-linux/recipes-core/bunch-update/bunch-update/update.service
[Unit]
Description=Bunch-Linux update service

[Service]
Type=oneshot
ExecStart=/usr/bin/bunch-update.sh

[Install]
WantedBy=multi-user.target
