# meta-bunch-linux/recipes-app/systemdservices/passivebuzzerservice.service
[Unit]
Description=Passive Buzzer service

[Service]
User=root
ExecStart=python3 /usr/bin/passiveBuzzer.py

[Install]
WantedBy=multi-user.target
